# bin/falzard.service
[Unit]
Description=Falzar daemon
After=network.target

[Service]
User=root
Group=root
EnvironmentFile=/etc/falzar/config.env
ExecStart=/usr/bin/falzard --port ${PORT} --data ${DATA}
RestartSec=5

[Install]
WantedBy=multi-user.target